# agent/pgagent.service.template
[Unit]
Description=PostgreSQL AI Agent
After=network.target postgresql.service

[Service]
Type=simple
User=pgagent
Group=pgagent
WorkingDirectory=/opt/pgagent/bin

# Variables injectées
Environment="AGENT_TOKEN=__TOKEN__"
Environment="AGENT_PORT=__PORT__"
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/opt/pgagent/bin"

# On lance via le venv
ExecStart=/opt/pgagent/venv/bin/python3 server.py

Restart=always
RestartSec=5

# Désactivation nécessaire pour Bubblewrap SUID
NoNewPrivileges=no

# Protection du système (lecture seule sur les dossiers vitaux)
ProtectSystem=full
# Empêche l'accès au dossier /home des autres utilisateurs
ProtectHome=true

[Install]
WantedBy=multi-user.target
